# etc/supervisor/conf.d/load_logger.conf
[program:load_logger]
command = /usr/bin/python /app/load_logger.py
autostart = true
autorestart = true
startsecs = 5
