# app/supervisor/logger.conf
[program:logger]
command = python /usr/src/app/logger.py
autostart = true
autorestart = true
startsecs = 5
