# app/supervisor/logger.conf
[program:logger]
command = python /app/logger.py
autostart = true
autorestart = true
startsecs = 5
